# ros/install_isolated/etc/catkin/profile.d/10.rosbuild.bat
REM generated from rosbuild/env-hooks/10.rosbuild.bat.em

REM env variables in installspace
set ROS_ROOT=F:/PROJECTS/infusion/ros-android/ros/install_isolated/share/ros
